# projects/tiktok-shorts/deploy/shorts-daily.timer
[Unit]
Description=Run tiktok-shorts every morning

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target

# 有効化:
#   cp deploy/shorts-daily.* /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now shorts-daily.timer
#   systemctl list-timers | grep shorts
